# infra/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server --save "" --appendonly no --loglevel notice
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:kcc-grpc]
command=/opt/venv/bin/python /opt/kcc-grpc/server.py
priority=20
autostart=true
autorestart=true
environment=KCC_SCRIPT="/opt/kcc/kcc-c2e.py",PYTHONPATH="/opt/kcc-grpc:/opt/kcc",GRPC_PORT="50051"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=/entrypoint-api.sh
priority=30
autostart=true
autorestart=true
environment=REDIS_URL="redis://localhost:6379",DATABASE_URL="%(ENV_DATABASE_URL)s",STORAGE_PATH="%(ENV_STORAGE_PATH)s",PORT="3000",NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=node /app/apps/worker/dist/main
directory=/app
priority=40
autostart=true
autorestart=true
environment=REDIS_URL="redis://localhost:6379",DATABASE_URL="%(ENV_DATABASE_URL)s",KCC_GRPC_ADDR="localhost:50051",NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
priority=50
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
